# apps/timemachine/startup-services.conf
[program:timemachine]
command=/opt/timemachine/start-timemachine-in-docker.sh
autorestart=true
