# indexers/notification-indexer/startup-services.conf
[program:logstash-indexer]
command=/opt/start-indexer.sh
autorestart=true
