# apps/search/startup-services.conf
[program:search-service]
command=/opt/cloudsight-kelk-search/start-search-service-in-docker.sh
autorestart=true
